# external-display/systemd/external-display@.service.template
[Unit]
Description=Scale UI on external display connection

[Service]
Type=oneshot
ExecStart={{ script_path }}

[Install]
WantedBy=graphical-session.target
