# external-display/udev/on-external-display-connection.rules.template
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", RUN+="{{ home }}/.local/bin/on-external-display-connection.sh"